# deployment/etf-trading.service
[Unit]
Description=ETF Trading System
After=network.target postgresql.service
Wants=postgresql.service

[Service]
Type=simple
User=etf-trader
Group=etf-trader
WorkingDirectory=C:\Users\hp\Desktop\ETF\ETF
Environment=PYTHONPATH=C:\Users\hp\Desktop\ETF\ETF
ExecStart=/usr/bin/python3 main.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=C:\Users\hp\Desktop\ETF\ETF/logs C:\Users\hp\Desktop\ETF\ETF/data

[Install]
WantedBy=multi-user.target
